# container-config/conf.d/pebble-emulator.conf
[program:pebble]
command=pebble install --emulator aplite --vnc
autorestart=true